# .ddev/supervisor/symfony-messenger-consume.conf
[program:shopware-consumer]
command=/usr/bin/php /var/www/html/shopware/bin/console messenger:consume async low_priority
autorestart=true
startretries=10
